# deploy/roboai-estun.service
[Unit]
Description=RoboAi Estun Codroid Driver
After=network.target

[Service]
Type=simple
User=teddy
Environment="ROS_DOMAIN_ID=0"
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/teddy/cobot_ws/install/setup.bash && ros2 run estun_driver estun_driver_node --ros-args -p robot_ip:=192.168.101.100 -p robot_port:=9000 -p poll_rate_hz:=50.0'
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
