# src/cobot_bringup/systemd/roboai-lidar.service
[Unit]
Description=RoboAi Livox MID-360 LiDAR
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teddy
WorkingDirectory=/home/teddy/cobot_ws
# LiDAR now lives on the isolated cell subnet (192.168.2.0/24) at
# 192.168.2.150, reached via eno1 (192.168.2.246). No /32 workaround
# needed — kernel routes 192.168.2.0/24 through eno1 natively. Do NOT
# reintroduce any 192.168.1.x address or route on eno1 here: the
# Jetson's WiFi (wlP1p1s0) owns 192.168.1.0/24 and any overlap on eno1
# hijacks the SSH return path.
ExecStartPre=+/sbin/ip link set eno1 up
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/teddy/cobot_ws/install/setup.bash && ros2 launch cobot_bringup lidar.launch.py'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
