# src/cobot_bringup/systemd/roboai-accumulator.service
[Unit]
Description=RoboAi Point Cloud Accumulator (near/far split)
After=roboai-lidar.service
Wants=roboai-lidar.service

[Service]
Type=simple
User=teddy
WorkingDirectory=/home/teddy/cobot_ws
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/teddy/cobot_ws/install/setup.bash && python3 /home/teddy/cobot_ws/src/cobot_bringup/scripts/pointcloud_accumulator.py'
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
