# src/cobot_bringup/systemd/roboai-lidar.service
[Unit]
Description=RoboAi Ouster LiDAR Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=teddy
WorkingDirectory=/home/teddy/cobot_ws
# Network setup runs as root via the '+' prefix so no sudoers entry is needed.
# We use a /32 + host route (not /24) because the Jetson WiFi is also on
# 192.168.1.0/24 — installing /24 on eno1 would hijack SSH return traffic.
ExecStartPre=+/sbin/ip link set eno1 up
ExecStartPre=+/bin/sh -c '/sbin/ip addr replace 192.168.1.200/32 dev eno1'
ExecStartPre=+/bin/sh -c '/sbin/ip route replace 192.168.1.100/32 dev eno1'
ExecStart=/bin/bash -c 'source /opt/ros/humble/setup.bash && source /home/teddy/cobot_ws/install/setup.bash && python3 /home/teddy/cobot_ws/src/cobot_bringup/scripts/ouster_bridge.py'
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
